# systemd/infra-maintenance.service
[Unit]
Description=Infra Maintenance Service

[Service]
ExecStart=/root/linux-infra-intern-assignment/scripts/maintenance.sh
